# tech.conf
[program:tech]
directory = /root/go/bin/
command= ./root/go/bin/eventica
autostart=true
autorestart=true
user = root
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log